# misc/pulsemeter/files/pulsemeter.in
#!/bin/sh
#
# PROVIDE: pulsemeter
# REQUIRE: DAEMON
# KEYWORD:

. /etc/rc.subr

name="pulsemeter"
rcvar="pulsemeter_enable"
pulsemeter_command="%%PREFIX%%/bin/pulsemeter"
pidfile="/var/run/${name}.pid"
command="/usr/sbin/daemon"
command_args="-P ${pidfile} -r -f ${pulsemeter_command}"

load_rc_config $name
: ${pulsemeter_enable:=no}

run_rc_command "$1"

